# deploy/aws/systemd/xiangqi-idle-check.timer
[Unit]
Description=Run the Xiangqi idle-shutdown check every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s

[Install]
WantedBy=timers.target
